# investor-robot.service
[Unit]
Description=investor — robot momentum multi-sector (orquestador 24/7)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/investor-app/momentum-investor
EnvironmentFile=/etc/investor.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/investor-app/venv/bin/python engine/orchestrator.py --loop
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
